# rouse-relay.service
[Unit]
Description=Rouse Relay Server
Documentation=https://oraserrata.net/relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rouse-relay
Restart=always
RestartSec=5

# Drop root. The relay only needs to bind a UDP socket on port 9876 and
# broadcast packets — neither of which requires privileges. DynamicUser
# allocates a transient user/group at runtime so we don't have to create
# one in /etc/passwd. Requires systemd 232+ (2016 and newer).
DynamicUser=true

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Environment — edit AUTH_TOKEN before enabling
Environment=AUTH_TOKEN=YOUR_PASSWORD_HERE
Environment=PORT=9876
Environment=HOST=0.0.0.0

# Logging goes to journald by default
# View with: journalctl -u rouse-relay -f

[Install]
WantedBy=multi-user.target
